# rootfs/etc/services.d/scheduler/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# 
# Home Assistant Add-on: SimpleScheduler
# 
# ==============================================================================

bashio::log.info "Starting service.d [Scheduler]"

exec /simplescheduler/scheduler.sh

