# rootfs/etc/services.d/apache2/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: SimpleScheduler
# Runs the Apache daemon
# ==============================================================================
bashio::log.info "Starting web server"
exec /usr/sbin/apache2ctl -D FOREGROUND
