# simplescheduler-beta/rootfs/etc/services.d/scheduler/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: SimpleScheduler
# Runs the scheduler daemon
# ==============================================================================

bashio::log.info "Starting scheduler daemon"

exec python3 /simplescheduler/scheduler.py
