# simplescheduler/rootfs/etc/services.d/interface/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: SimpleScheduler
# Runs the inteface daemon
# ==============================================================================

bashio::log.info "Starting main daemon"

exec python3 /simplescheduler/main.py
